# scripts/bandwidth-guard.service
[Unit]
Description=Bandwidth Guard - Network Usage Monitor
After=network.target
Documentation=https://github.com/Ndukiye/bandwidth-guard

[Service]
Type=simple
User=root
ExecStart=/opt/bandwidth-guard/venv/bin/python /opt/bandwidth-guard/src/monitor.py
WorkingDirectory=/opt/bandwidth-guard
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Desktop notification support
Environment="DISPLAY=:0"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus"

[Install]
WantedBy=multi-user.target